# rootfs/etc/periodic/hourly/spotweb
#!/bin/sh

s6-setuidgid abc php /app/spotweb/retrieve.php
